# esp32-auto-program@.service
[Unit]
Description=Program ESP32 device attached to '%I'

[Service]
User=<user>
ExecStart=<path>/esp32-program-shim.sh '%I'
